# templates/systemd/guestbook.service.j2
[Unit]
Description=guestbook service
[Service]
WorkingDirectory={{project_path}}
ExecStart=/usr/bin/python3 {{project_path}}
Type=simple